# loyverse-api.service
[Unit]
Description=Loyverse Automation API
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/path/to/your/project
Environment=DISPLAY=:99
Environment=NODE_ENV=production
ExecStartPre=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
ExecStart=/usr/bin/pm2 start ecosystem.config.js --env production
ExecReload=/usr/bin/pm2 reload loyverse-api
ExecStop=/usr/bin/pm2 stop loyverse-api
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
